# etc/supervisor/conf.d/dvv.conf
[program:vnc]
command=bash -c '/dockerstartup/go_go.sh'
autostart=true
autorestart=true
loglevel=info

