# etc/supervisor/conf.d/sshd.conf
[supervisord]
nodaemon=true
[program:SSH]
command=/usr/sbin/sshd -D -e
stderr_logfile=/var/log/myapi.err.log
stdout_logfile=/var/log/myapi.out.log
directory=/root/

